# deploy/grab-a1.service
[Unit]
Description=portal-mapper A1 upgrader (runs on the 1-OCPU A1, resizes it to 2 OCPU 24/7)
Documentation=https://github.com/your/repo/blob/master/deploy/README.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Adjust path/user to wherever you put the repo on the A1 box.
ExecStart=/home/opc/portal-mapper/deploy/grab-a1.sh
# Optional: phone push when the box reaches 2 OCPU (https://ntfy.sh, no signup).
# Environment=NTFY_TOPIC=your-secret-topic

# The script loops forever internally and exits only on terminal outcomes:
#   0 = resize accepted / already at 2 OCPU   2 = auth broken   3 = limit hit
# A successful resize reboots the box; the unit simply re-runs on boot, sees the
# box is already at 2 OCPU, and exits 0 again (idempotent). Treat 0/2/3 as
# "done, do NOT restart on exit"; restart only if it dies unexpectedly (network
# drop, killed) so it self-heals across blips and reboots.
SuccessExitStatus=0 2 3
RestartPreventExitStatus=2 3
Restart=on-failure
RestartSec=15

Nice=10

[Install]
WantedBy=multi-user.target
